# applications/insomnia/root/etc/services.d/x11vnc/run
#!/usr/bin/execlineb -P

set -eu

s6-svwait /var/run/s6/services/xvfb

x11vnc -storepasswd 'toor' /etc/x11vnc.pass
x11vnc -rfbauth /etc/x11vnc.pass -display :1 -forever
