# applications/insomnia/root/etc/services.d/x11vnc/run
#!/usr/bin/execlineb -P
x11vnc -display :1 -forever
