# applications/insomnia/root/etc/services.d/novnc/run
#!/usr/bin/execlineb -P
/opt/novnc/utils/novnc_proxy --vnc localhost:5900
