# applications/insomnia/root/etc/services.d/insomnia/run
#!/usr/bin/execlineb -P
with-contenv
DISPLAY=:1 insomnia
